# fraisier/scaffold/templates/core/fraisier-webhook.service.j2
[Unit]
Description=Fraisier webhook server for {{ project_name }}
After=network.target
Wants=network-online.target

[Service]
Type=exec
User={{ scaffold.deploy_user }}
Group={{ scaffold.deploy_user }}
ExecStart=/home/{{ scaffold.deploy_user }}/.local/bin/fraisier-webhook
Restart=on-failure
RestartSec=5
# WatchdogSec intentionally omitted: uvicorn does not implement sd_notify
# watchdog pings, so WatchdogSec would kill the process every N seconds.

# Environment
EnvironmentFile=/etc/fraisier/{{ project_name }}.webhook.env
{% if pg_allowed_databases %}
Environment=FRAISIER_PG_WRAPPER=/usr/local/libexec/fraisier/pgadmin-{{ project_name }}
{% endif %}
{% if allowed_services %}
Environment=FRAISIER_SYSTEMCTL_WRAPPER=/usr/local/libexec/fraisier/systemctl-{{ project_name }}
{% endif %}

# Resource limits
MemoryMax=512M
CPUAccounting=true

# Security hardening
# ProtectHome intentionally omitted: fraisier-webhook is installed as a uv tool
# and runs from ~/.local, which ProtectHome (any value) makes inaccessible.
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Runtime directory (created by systemd at startup, owned by deploy user)
RuntimeDirectory=fraisier
RuntimeDirectoryMode=0750

# Writable paths
ReadWritePaths=/var/lib/fraisier
ReadWritePaths=/run/fraisier
ReadWritePaths=/opt/fraisier
{% for fraise in local_fraises %}
{% for env_name, env_config in fraise.get('environments', {}).items() %}
{% if env_config.get('app_path') %}
ReadWritePaths={{ env_config['app_path'] }}
{% endif %}
{% endfor %}
{% endfor %}

[Install]
WantedBy=multi-user.target
