# fraisier/scaffold/templates/core/poll-deploy.service.j2
[Unit]
Description=Poll-deploy checker for {{ project_name }}
After=network.target

[Service]
Type=oneshot
User={{ scaffold.deploy_user }}
{% if pg_allowed_databases %}
Environment=FRAISIER_PG_WRAPPER=/usr/local/libexec/fraisier/pgadmin-{{ project_name }}
{% endif %}
{% if allowed_services %}
Environment=FRAISIER_SYSTEMCTL_WRAPPER=/usr/local/libexec/fraisier/systemctl-{{ project_name }}
{% endif %}
ExecStart=/usr/local/bin/fraisier trigger-deploy {{ fraise_names | first }} production --force

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
ReadWritePaths=/var/lib/fraisier /var/log/fraisier
{% for fraise in fraises %}
ReadWritePaths=/opt/{{ fraise.name }}
{% endfor %}
